# scripts/com.vibox.prune.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vibox.prune</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/vimo/vimo-cloud/scripts/vibox-prune.sh</string>
    </array>

    <key>UserName</key>
    <string>vimo</string>

    <!-- 매일 새벽 3시 (저트래픽 시간대) -->
    <key>StartCalendarInterval</key>
    <dict>
      <key>Hour</key><integer>3</integer>
      <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.vibox.prune.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.vibox.prune.err</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
  </dict>
</plist>
